# vectorbot2-backup.timer
[Unit]
Description = Backup vectorbot2 database

[Timer]
OnCalendar = weekly

[Install]
WantedBy = timers.target
